# deployment/systemd/kucoin-basis-convergence-strategy.service
[Unit]
Description=KuCoin basis-convergence paper strategy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kucoin-basis-funding-arb
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/kucoin-basis-funding-arb/.venv/bin/python kucoin_basis_convergence/run_paper_strategy.py --loop --interval 60
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
